# ansible/haproxyconfd/confd.service
[Unit]
Description=Confd
After=haproxy.service

[Service]
ExecStart=/usr/sbin/confd
Restart=always

[Install]
WantedBy=basic.target
